# .docker/services/esplora/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
username=dummy
password=dummy

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:electrs]
directory=/usr/local/src/electrs
command=/etc/init.d/electrs.sh start
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:esplora]
directory=/usr/local/src/esplora
environment=API_URL="http://localhost:3000/"
command=/etc/init.d/esplora.sh start
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=dummy
password=dummy